# src/root/etc/systemd/system/post-run.service
[Unit]
Description=Pire psot-run stuff

[Service]
ExecStart=/etc/post-run.sh

[Install]
WantedBy=multi-user.target